# templates/puppetdb-DLO-cleanup.service.epp
<%- | String[1] $puppetdb_user,
      String[1] $puppetdb_group,
      String[1] $vardir,
      Integer[1] $dlo_max_age
| -%>
[Unit]
Description=Cleanup old discarded puppetdb reports

[Service]
Type=oneshot
User=<%= $puppetdb_user %>
Group=<%= $puppetdb_group %>
ExecStart=/usr/bin/find <%= $vardir %>/stockpile/discard/ -type f -mtime +<%= $dlo_max_age %> -delete
